# clients/crontab
00 01 * * * root /usr/local/scripts/rsyncup > /dev/null 2> /dev/null